# qphone/com.qphone.webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.qphone.webhook</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>/Users/q/agents/qagents/qphone/index.js</string>
		<string>--server</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/q/.openclaw/workspace/logs/qphone-webhook.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/q/.openclaw/workspace/logs/qphone-webhook.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>QPHONE_WEBHOOK_URL</key>
		<string>REPLACE_WITH_CLOUDFLARE_TUNNEL_URL</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/q/agents/qagents/qphone</string>
</dict>
</plist>
